# systemd/goodwe-stack.target
[Unit]
Description=GoodWe control stack (control + ingest + api + ui)
Wants=goodwe-control.service goodwe-ingest.service goodwe-api.service goodwe-ui.service
